# setup/codeScraperTimer.timer
[Unit]
Description=Runs the PheonixCodeScraper file once every n (7) days

[Timer]
OnCalendar=Mon *-*-* 00:00:00

[Install]
WantedBy=timers.target